# packages/solidus/solidus-0.10.5.tar.gz/solidus-0.10.5/solidus/solidus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solidus.report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Python/2.7/site-packages/solidus-0.10.00-py2.7.egg/solidus/Solidus.py</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
</dict>
</plist>
