# packages/btcndash/btcndash-2.3.6.tar.gz/btcndash-2.3.6/scripts/bitcoind.conf
########################################
# Upstart script for Bitcoin Full Node
# Save as /etc/init/bitcoin.conf
# then do:
#   > sudo initctl reload-configuration
#   > sudo start bitcoind
########################################

description "Bitcoin Full Node"

# Set the location of the log file
env LOG_FILE=/home/bitcoin/log/bitdash.log

# Set the location of the bitdash python file
env DAEMON=/home/bitcoin/bin/32/bitcoind

# Set the user and group to use here (don't use root!)
setuid bitcoin
setgid bitcoin

# Sets when this script will run
start on runlevel [2345]
stop on runlevel [!2345]

# Start the daemon
exec $DAEMON >> $LOG_FILE 2>&1

# Respawn if crashed, but only a few times
respawn
respawn limit 10 5
