# packages/btcndash/btcndash-2.3.6.tar.gz/btcndash-2.3.6/scripts/btcndash.conf
########################################
# Upstart script for BTCnDash
# Save as /etc/init/btcndash.conf
# then do:
#   > sudo initctl reload-configuration
#   > sudo start btcndash
########################################

description "BTCnDash Bitcoin Node Status Dashboard"

# Set the root directory for your virtualenv here
env PYTHON_HOME=/home/btcndash/venv

# Set the location of the log file
env LOG_FILE=/home/btcndash/log/btcndash.log

# Set the location of the bitdash python file
env DAEMON=/home/btcndash/venv/lib/python2.7/site-packages/btcndash/btcndash.py
env OPTIONS=--config=/home/btcndash/config.json

# Set the user and group to use here (don't use root!)
setuid btcndash
setgid btcndash

# Sets when this script will run
start on runlevel [2345]
stop on runlevel [!2345]

# Start the daemon
exec $PYTHON_HOME/bin/python $DAEMON $OPTIONS >> $LOG_FILE 2>&1

# Respawn if crashed, but only a few times
respawn
respawn limit 10 5

