# packages/portinus/portinus-0.9.4.tar.gz/portinus-0.9.4/portinus/templates/instance.service
[Unit]
Description={{ name }} service
After=network.target docker.service

[Service]
Type=simple
{% if environment_file -%}
EnvironmentFile={{ environment_file.path }}
{% endif -%}
ExecStart={{ start_command }}
ExecStop={{ stop_command }}

[Install]
WantedBy=multi-user.target
